# resources/qtile-session.target
[Unit]
Description=qtile session
BindsTo=graphical-session.target
Before=graphical-session.target
